# roles/data-scraper/files/rc.local
#!/bin/bash
/bin/sh -cf  '/usr/bin/svscanboot &'
exit 0
